# systemd/user@.service.d/10-rpi-preseed.conf
# Ordering-only: hold every user session until base config + early runcmds are done.
# The late phase is deliberately excluded so logins/SSH/desktop stay available
# during its (potentially long) network/apt waits. Never Requires= — a preseed
# failure or later-boot skip must not lock users out.
[Unit]
After=rpi-preseed.service rpi-preseed-runcmd-early.service
